# autostart/google-chrome.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/google-chrome
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en]=Chrome
Name=Chrome
Comment[en]=
Comment=
